# libreoffice-cups/etc/systemd/user/cups-fuse.path
[Unit]
Description=Starts mounting printers to %h/Desktop/Drucker, if this directory exists
Wants=set-environment.service
After=set-environment.service

[Path]
PathExists=%h/Drucker
Unit=cups-fuse.service

[Install]
WantedBy=paths.target
